# backend/supervisord.conf
[supervisord]
nodaemon=true
logfile=/app/logs/supervisord.log
pidfile=/var/run/supervisord.pid
user=root

[program:redis]
command=redis-server --save 60 1 --loglevel warning --dir /app/data
autostart=true
autorestart=true
stdout_logfile=/app/logs/redis.log
stderr_logfile=/app/logs/redis_err.log
priority=1
startsecs=5

[program:celery-worker]
command=celery -A celery_app worker --loglevel=info --concurrency=4
directory=/app
autostart=true
autorestart=true
stdout_logfile=/app/logs/celery.log
stderr_logfile=/app/logs/celery_err.log
priority=2
startsecs=10
environment=REDIS_HOST="localhost",REDIS_PORT="6379",REDIS_DB="0",CONFIG_DIR="/app/data"

[program:flask-api]
command=python3 web_api.py --host 0.0.0.0 --port 5000
directory=/app
autostart=true
autorestart=true
stdout_logfile=/app/logs/flask.log
stderr_logfile=/app/logs/flask_err.log
priority=3
startsecs=10
environment=REDIS_HOST="localhost",REDIS_PORT="6379",REDIS_DB="0",CONFIG_DIR="/app/data"
